# farmprices/deploy/farmprices.service
[Unit]
Description=Tenbury Farm Supplies Price Lookup App
After=network.target
Wants=network-online.target

[Service]
Type=exec
# Change 'richowen' if your Pi username is different
User=richowen
Group=richowen
WorkingDirectory=/home/richowen/Inventory-Prices/farmprices
Environment="PATH=/home/richowen/Inventory-Prices/farmprices/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=-/home/richowen/Inventory-Prices/farmprices/.env

# DB init happens inside create_app() — no ExecStartPre needed
ExecStart=/home/richowen/Inventory-Prices/farmprices/venv/bin/gunicorn \
    --workers 4 \
    --bind 0.0.0.0:5000 \
    --timeout 120 \
    --access-logfile /home/richowen/Inventory-Prices/farmprices/app.log \
    --error-logfile /home/richowen/Inventory-Prices/farmprices/app.log \
    --log-level info \
    "app:create_app()"

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=farmprices

[Install]
WantedBy=multi-user.target
